# systemd/lmi-sync@.service
[Unit]
Description=Run sync-all once for lmi-service
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/jose/%i
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/jose/%i/.venv/bin/python scripts/sync_all.py
User=jose
Group=jose
