# systemd/lmi-sync@.service
[Unit]
Description=Run sync-all once for lmi-service
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/%i
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/%i/.venv/bin/python scripts/sync_all.py
User=jose
Group=jose
