# systemd/lmi@.service
[Unit]
Description=LunchMoney Ingest Service (lmi-service)
After=network.target

[Service]
WorkingDirectory=%h/%i
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/%i/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8010
Restart=on-failure
User=jose
Group=jose

[Install]
WantedBy=multi-user.target
